# postgresql-source-debug/config/supervisord.conf
[supervisord]
nodaemon=true

[program:postgresql]
user=postgres
command=postgres -D /var/lib/pgsql/13/data -c config_file=/etc/postgresql/postgresql.conf